# cam/cam.timer
[Unit]
Description=Capture a snapshot hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
