# cam/cam.service
[Unit]
Description=Camera snapshot

[Service]
ExecStart=/root/cam/s3snapshot.sh
